# meta-intel-quark/recipes-core/images/files/rfkill-unblock.service
[Unit]
Description=unblock script
After=systemd-rfkill@.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/rfkill unblock all
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

