# meta-tools/recipes-test/stress-script/files/power-test.service
[Unit]
Description=Power Test

[Service]
Type=oneshot
ExecStart=/tools/stress.sh

[Install]
WantedBy=multi-user.target
