# meta-tools/recipes-test/emi-script/files/emi-test.service
[Unit]
Description=EMI Test

[Service]
Type=oneshot
ExecStart=/tools/emi_run

[Install]
WantedBy=multi-user.target
