# meta-tools/recipes-test/boottimes-script/files/boottimes.service
[Unit]
Description=Boot Times

[Service]
Type=oneshot
ExecStart=/tools/boottimes.sh
StandardOutput=kmsg+console

[Install]
WantedBy=multi-user.target
